# defi_stats.service
[Unit]
Description=fastapi_defi_stats
Wants=network-online.target
After=network-online.target
Conflicts=getty@tty1.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Environment=PATH=/home/admin/.local/bin:$PATH
WorkingDirectory=/home/admin/defi_stats
ExecStart=/home/admin/defi_stats/run_api.sh

# Logs path
StandardOutput=append:/home/admin/logs/defi-stats.log
StandardError=append:/home/admin/logs/defi-stats.log

User=admin
Group=admin
Type=simple
TimeoutStopSec=30min
Restart=on-failure
RestartSec=10s
StandardInput=tty-force

[Install]
WantedBy=multi-user.target

